# vpn-server/systemd/vpn-backup.timer
[Unit]
Description=vpn-server: ежедневный бэкап конфигов

[Timer]
OnCalendar=daily
RandomizedDelaySec=30m
Persistent=true

[Install]
WantedBy=timers.target
