# vpn-server/systemd/vpn-backup.service
[Unit]
Description=vpn-server: бэкап конфигов wg-easy/Xray/routing
Wants=network-online.target
After=network-online.target docker.service

[Service]
Type=oneshot
WorkingDirectory=__VPN_ROOT__
ExecStart=__VPN_ROOT__/scripts/backup.sh
